# scripts/systemd/system/pigpiod.service.d/public.conf
[Unit]
Description=pigpio daemon

[Service]
ExecStart=
ExecStart=/usr/local/bin/pigpiod -k
PIDFile=/var/run/pigpio.pid
User=root
Restart=always
RestartSec=5s
